# rules/windows/process_creation/scheduled_task_persistence.yml
title: Scheduled Task Creation via schtasks
id: 2f7d9a3e-4c81-4b25-9e0f-7a1c6d8b2e94
status: experimental
description: |
  Detects creation of a scheduled task from the command line via schtasks.exe.
  Adversaries abuse scheduled tasks for persistence and privileged execution,
  frequently pointing the action at a script interpreter or a file in a
  user-writable directory.
references:
  - https://attack.mitre.org/techniques/T1053/005/
author: Vandit Shah (BLUE-FORGE)
date: 2026-07-12
tags:
  - attack.persistence
  - attack.privilege-escalation
  - attack.execution
  - attack.t1053.005
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
    - Image|endswith: '\schtasks.exe'
    - OriginalFileName: 'schtasks.exe'
  selection_create:
    CommandLine|contains: '/create'
  selection_suspicious:
    CommandLine|contains:
      - '\AppData\'
      - '\Temp\'
      - '\Users\Public\'
      - 'powershell'
      - 'cmd /c'
      - 'mshta'
      - 'rundll32'
      - '.vbs'
      - '.js'
  condition: selection_img and selection_create and selection_suspicious
falsepositives:
  - Software installers and management tools that register scheduled tasks.
    Baseline task names and calling processes in your environment.
level: medium
